# bo7metrv.service
[Unit]
Description=bo7metrv image and servo controller
After=local-fs.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bo7metrv
ExecStartPre=/bin/sleep 8
ExecStart=/usr/bin/python3 /home/pi/bo7metrv/run.py
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target

# The touchscreen uses SPI0 plus PENIRQ on GPIO25.
# Add these boot config lines to /boot/firmware/config.txt or /boot/config.txt:
# dtparam=spi=on
# dtoverlay=ads7846,penirq=25,cs=1
#
# Install/update the service with:
# sudo cp bo7metrv.service /etc/systemd/system/bo7metrv.service
# sudo systemctl daemon-reload
# sudo systemctl enable --now bo7metrv.service
